# deploy/bobrik-miniapp.service
[Unit]
Description=Bobrik Mini App FastAPI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/my-private
Environment=PYTHONUNBUFFERED=1
Environment=BOT_TOKEN=replace_me
ExecStart=/usr/bin/python3 -m uvicorn webapp.backend.main:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
